# kafka/supervisor/kafka.conf
[program:kafka]
command=/usr/bin/start-kafka.sh
autostart=true
autorestart=true

;Print Kafka output to console, use with start-foreground
;stdout_logfile=/dev/fd/1
;stdout_logfile_maxbytes=0